# gammastep/.config/systemd/user/gammastep.service
[Unit]
Description=gammastep

[Service]
Type=simple
StandardOutput=journal
ExecStart=/usr/bin/gammastep

[Install]
WantedBy=graphical.target
